# zabbix-doc/zabbix/etc/cron.d/nc-zabbix
##############################################
## NC zabbix cronjob for application
##############################################
*/2 * * * * nc-zabbix bash /var/lib/nc_zabbix/cron/nc_iostat_cron.sh
*/2 * * * * nc-zabbix bash /var/lib/nc_zabbix/cron/nc_iptables_cron.sh
*/2 * * * * nc-zabbix bash /var/lib/nc_zabbix/cron/nc_vmstat_cron.sh
*/2 * * * * nc-zabbix bash /var/lib/nc_zabbix/cron/nc_mpstat_cron.sh
#0 */12 * * * nc-zabbix bash /var/lib/nc_zabbix/cron/nc_megacli_cron.sh
#*/5 * * * * nc-zabbix bash /var/lib/nc_zabbix/cron/nc_xen_cron.sh

##############################
# nc-zabbix-agent reload script#
##############################
#*/5 * * * * root bash /var/lib/nc_zabbix/zabbix_agent_reloader.sh
